# systemd/microsoft-job-watcher.service
[Unit]
Description=Microsoft Careers software engineer watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_LINUX_USER
Group=YOUR_LINUX_USER
WorkingDirectory=/opt/microsoft-job-watcher
EnvironmentFile=/opt/microsoft-job-watcher/microsoft-job-watcher.env
ExecStart=/usr/bin/python3 /opt/microsoft-job-watcher/microsoft_job_watcher.py --interval-minutes 15 --max-pages 200 --stop-after-seen-pages 3 --full-scan-interval-hours 24 --display-timezone America/Los_Angeles --webhook-url http://127.0.0.1:18789/hooks/agent --webhook-mode openclaw-agent --webhook-bearer-token-env OPENCLAW_HOOKS_TOKEN --openclaw-agent-id microsoft-jobs --openclaw-session-key hook:microsoft-jobs --openclaw-channel telegram --openclaw-to YOUR_TELEGRAM_USER_ID --no-print-matches
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
